# hassio-hostapd/rootfs/etc/services.d/ap.svc/run
#!/usr/bin/execlineb -P
# ==============================================================================
# Home Assistant Community Add-on: Example
# Runs example1 script
# ==============================================================================


#bashio::log.info "Starting Main..."

exec /usr/bin/ap.sh